# atd.service.in
[Unit]
Description=Deferred execution scheduler

[Service]
ExecStart=@sbindir@/atd -f

[Install]
WantedBy=multi-user.target
